# deploy/systemd/musicseer.service
[Unit]
Description=MusicSeer Enhanced — music discovery and requests for Lidarr
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=musicseer
Group=musicseer
ExecStart=/opt/musicseer/musicseer
WorkingDirectory=/opt/musicseer
Environment=MUSICSEER_DATA_DIR=/var/lib/musicseer
Environment=MUSICSEER_PORT=8688
EnvironmentFile=-/etc/musicseer/musicseer.env
Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
# /opt/musicseer is writable so the in-app updater can replace the binary and
# re-exec. It must also be owned by the musicseer user (chown -R). Downloads are
# SHA-256-verified against the release before the swap. Drop /opt/musicseer here
# (and from the chown) if you prefer CLI-only upgrades.
ReadWritePaths=/var/lib/musicseer /opt/musicseer
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
